# recipes-apps/server_app/files/systemd-start-server.service
[Unit]
Description=Service file to Start server

[Service]
ExecStart=/usr/bin/server
Restart=always
User=root

[Install]
WantedBy=multi-user.target
